# init/mapi
#!/bin/ksh

daemon="/home/project-wildcard/backend/server.js"

. /etc/rc.d/rc.subr

rc_start() {
        ${rcexec} "${daemon} ${daemon_flags} 2>&1 | logger -t mapi &"
}

rc_cmd $1